# templates/backup_dns.service.j2
[Unit]
Description=Backup DNS service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/bash /usr/local/bin/backup_dns.sh
