# templates/backup_dns.timer.j2
[Unit]
Description=Run backup job every 5 minutes

[Timer]
OnUnitActiveSec=5m
Unit=backup_dns.service

[Install]
WantedBy=timers.target
